# cron/crontab
# Run indexer every hour at minute 5
5 * * * * docker exec fs-indexer python -u /app/indexer.py 2>&1